# .config/profile.d/bkt
# BKT cache command outputs

export BKT_CACHE_DIR="$XDG_CACHE_HOME/bkt"

