# .config/profile.d/cargo
# Cargo
export CARGO_INSTALL_ROOT=$HOME/.local
export CARGO_HOME=$XDG_CONFIG_HOME/cargo

