# .config/profile.d/composer
# composer
export COMPOSER_BIN_DIR="$HOME/.local/bin"
export COMPOSER_CACHE_DIR="$XDG_CACHE_HOME/composer"
export COMPOSER_HOME="$XDG_CONFIG_HOME/composer"

