# .config/profile.d/npm
# npm
export NPM_CONFIG_USERCONFIG="$XDG_CONFIG_HOME/npm.conf"
export NPM_CONFIG_CACHE="$XDG_CACHE_HOME/npm"
export NPM_CONFIG_PREFIX="$HOME/.local"

