# .config/profile.d/ssh
if ! pidof ssh-agent ; then
  eval $(ssh-agent -s)
fi

